# ansible/roles/postgresql/templates/pg-backup.timer.j2
# ── PostgreSQL Daily Backup Timer (systemd timer) ────────────────────
# Triggers pg-backup.service daily at 02:00 UTC.
[Unit]
Description=Lazynext PostgreSQL Daily Backup Timer
Requires=pg-backup.service

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=600
AccuracySec=5min

[Install]
WantedBy=timers.target
