# ansible/roles/postgresql/templates/pg-backup.service.j2
# ── PostgreSQL Daily Backup Service (systemd unit) ───────────────────
# Runs pg_dump with gzip compression and stores locally.
# Triggered by pg-backup.timer at 02:00 UTC daily.
[Unit]
Description=Lazynext PostgreSQL Daily Backup
After=network.target postgresql@{{ postgresql_version }}-main.service
Requires=postgresql@{{ postgresql_version }}-main.service

[Service]
Type=oneshot
User=postgres
Group=postgres
ExecStart=/usr/local/bin/pg_backup.sh
Environment=PGPASSWORD={{ postgresql_admin_password }}
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pg-backup
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ postgresql_backup_dir }} /tmp
